# apps/whoami/whoami.container
[Unit]
Description=whoami
After=network-online.target
Wants=network-online.target

[Container]
AddCapability=NET_BIND_SERVICE
ContainerName=whoami
Environment=TZ=Europe/London
Image=docker.io/traefik/whoami:latest
Label="traefik.enable='true'"
Label="traefik.http.routers.whoami.rule=Host(`traefik.kvshs.xyz/whoami`)"
Label="traefik.http.routers.whoami.tls=true"
Network=host

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
